# curator4-1452402603.conf
[program:curator4-1452402603]
command=python3 /home/curator4/curator4-1452402603/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1452402603/test.err.log
stdout_logfile=/home/curator4/curator4-1452402603/test.out.log